# celery-beat-rocktea.service
[Unit]
Description=Celery Beat Scheduler for Rocktea Mall Django project
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/django-app/dev/main
EnvironmentFile=/home/ubuntu/django-app/dev/main/setup/.env

# Beat scheduler configuration
ExecStart=/home/ubuntu/django-app/dev/venv/bin/celery -A setup beat \
    --loglevel=info \
    --schedule=/var/run/celery/celerybeat-schedule \
    --pidfile=/var/run/celery/beat.pid

# Restart configuration
Restart=always
RestartSec=10s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=append:/var/log/celery/beat.log
StandardError=append:/var/log/celery/beat.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target